# voice-journal-debrief.service
[Unit]
Description=Voice Journal — Morning Debrief Email
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=journal
Group=journal
WorkingDirectory=/opt/voice-journal
# NOTE: do NOT add EnvironmentFile=.env — the app loads .env via python-dotenv.
# systemd's EnvironmentFile parser is stricter and fails on values dotenv accepts
# (spaces, special chars), which makes the service fail before Python runs.
ExecStart=/opt/voice-journal/venv/bin/python3 /opt/voice-journal/debrief/main.py
StandardError=append:/opt/voice-journal/debrief/logs/debrief.log
