# systemd/voice-journal-healthcheck.timer
[Unit]
Description=Voice Journal Healthcheck — 08:00 daily

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
